# config-files/systemd/dcgm.service.in
[Unit]
Description=DEPRECATED. Please use nvidia-dcgm.service
Conflicts=nvidia-dcgm.service

[Service]
User=root
PrivateTmp=false

ExecStart=@DCGM_HOSTENGINE_LOCATION@/nv-hostengine -n --service-account nvidia-dcgm

Restart=on-abort

[Install]
WantedBy=multi-user.target
