# config-files/systemd/nvidia-dcgm.service.in
[Unit]
Description=NVIDIA DCGM service
Conflicts=dcgm.service

[Service]
User=root
PrivateTmp=false

ExecStart=@DCGM_HOSTENGINE_LOCATION@/nv-hostengine -n

Restart=on-abort

[Install]
WantedBy=multi-user.target
